# scripts/systemd/atw-web-backend.service
[Unit]
Description=ATW Web UI - Backend (FastAPI)
After=network.target

[Service]
Type=simple
User=%USER%
Group=%USER%
WorkingDirectory=%PROJECT_ROOT%/backend
Environment="PATH=%PROJECT_ROOT%/backend/.venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=%PROJECT_ROOT%/backend/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

# Logging
StandardOutput=append:%PROJECT_ROOT%/logs/backend.log
StandardError=append:%PROJECT_ROOT%/logs/backend.log

[Install]
WantedBy=multi-user.target
